# newmake.files
cdk/Makefile.am
cdk/configure.ac
cdk/acinclude.m4
cdk/doc/BUILD.changes
cdk/Patches/busybox.config.m4
cdk/Patches/dbox2-flash.c.m4
cdk/root/Makefile.inc
cdk/root/etc/init.d/Makefile.am
cdk/root/etc/init.d/rcS.m4
cdk/root/share/udhcpc/Makefile.am
cdk/backwards_compat.mk
cdk/root/etc/Makefile.am
apps/tuxbox/plugins/configure.ac
apps/tuxbox/tools/configure.ac
apps/tuxbox/tools/misc/Makefile.am
